# supervisor.ini
; this is the configuration for supervisor
; "fab webserver deploy" will copy it to /etc/supervisor/project_name.ini

[program:project_name]
command=/var/www/project_name/bin/python /var/www/project_name/releases/current/project_name/manage.py run_gunicorn -c /var/www/project_name/releases/current/gunicorn-settings.py
process_name=project_name-gunicorn
autorestart=true
user=project_name
stdout_logfile=/var/www/project_name/logs/supervisor.log
stdout_logfile_maxbytes=1MB
;stdout_logfile_backups=10
stderr_logfile=/var/www/project_name/logs/supervisor_error.log
stderr_logfile_maxbytes=1MB
;stderr_logfile_backups=10
directory=/var/www/project_name
